# resources/contrib/crontab
# Example crontab for PIVX Name Indexer
# Adjust paths and user permissions according to your deployment.

# Run every 1 minute
* * * * * /usr/bin/php /var/www/pivx-name-indexer/indexer.php sync >/dev/null 2>&1
