# resources/contrib/pivx-name-indexer.service
# Adjust paths and user permissions according to your deployment.

[Unit]
Description=PIVX Name Indexer Sync Service
After=network.target mariadb.service mysql.service rr-proxy-indexer.service

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/pivx-name-indexer
ExecStart=/usr/bin/php indexer.php sync

[Install]
WantedBy=multi-user.target
